# crontab
0 3 1 * * root /usr/local/bin/python /app/run.py >> /var/log/cron.log 2>&1
